# rpm/after-install
#!/bin/sh

/sbin/chkconfig --add twemproxy
